# debian/beaver.upstart
description "Beaver"

start on runlevel [2345]
stop on runlevel [!2345]
respawn

env BEAVER_MODE=connect

exec /usr/bin/beaver -c /etc/beaver.conf
